# lemonitor.timer
[Unit]
Description=Invoke lemonitor daily

[Timer]
OnBootSec=10min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
